# assets/vpncloud.target
[Unit]
Description=VpnCloud target allowing to start/stop all vpncloud@.service instances at once
